# src/init.d/disks-spindown
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

supervisor="supervise-daemon"
command="/usr/bin/disks-spindown /etc/conf.d/disks-spindown"
description="Daemon for Advanced Configuration and Power Interface"

depend() {
	use logger
}
